# deployment/systemd/fci-chatbot.service
[Unit]
Description=AI Chatbot Service
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/ai-chatbot
ExecStart=/usr/bin/docker-compose up -d
ExecStop=/usr/bin/docker-compose down
TimeoutStartSec=0

# Security settings
User=root
Group=root
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/ai-chatbot /opt/ai-chatbot/logs

# Restart policy
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target 